# ansible/templates/movie-api.service.j2
[Unit]
Description=Movie Analyst API
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/movie-analyst-api
ExecStart=/usr/bin/npm start
Restart=always
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
